# scripts/launchd/com.peachtreetowers.pricing-23n.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.peachtreetowers.pricing-23n</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/yasserkhan/airbnb-autohost2/scripts/daily-pricing-23n.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/yasserkhan/airbnb-autohost2</string>

    <!-- Daily at 09:00 local time -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Don't fire on load; only on the 09:00 schedule -->
    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/yasserkhan/airbnb-autohost2/logs/launchd.23n.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/yasserkhan/airbnb-autohost2/logs/launchd.23n.err.log</string>
</dict>
</plist>
